# mcro.service
[Unit]
Description=mcro - macros and shortcuts

[Service]
ExecStart=%h/.cargo/bin/mcro
Restart=on-failure
RestartSec=1
