# mcro.service
[Unit]
Description=mcro - macros and shortcuts

[Service]
Environment="PATH=%h/.cargo/bin"
ExecStart=mcro
